# additional_essential_commands/touchfile.service
[Unit]
Description=Touch File Service

[Service]
Type=oneshot
ExecStart=/usr/bin/touch /tmp/touched_file
