# additional_essential_commands/web.socket
[Unit]
Description=Python HTTP Server Socket
[Socket]
ListenStream=8080
[Install]
WantedBy=sockets.target