# storage/data-mount.automount
# /etc/systemd/system/data-mount.automount
[Unit]
Description=Automount My Mount

[Automount]
Where=/data/mount
TimeoutIdleSec="10 min"

[Install]
WantedBy=multi-user.target